# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/callis/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid
loglevel=info

[program:api]
command=/app/.venv/bin/python main.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=APP_VERSION="%(ENV_APP_VERSION)s"

[program:sshd]
command=/entrypoint-sshd.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
